# systemd/system/snap-core-2312.mount
[Unit]
Description=Mount unit for core

[Mount]
What=/var/lib/snapd/snaps/core_2312.snap
Where=/snap/core/2312
Type=squashfs
Options=nodev,ro

[Install]
WantedBy=multi-user.target
